# deploy/systemd/wireportal-api.service
[Unit]
Description=WirePortal API service
Wants=network-online.target
After=network-online.target wireportal-wg-agent.service
Requires=wireportal-wg-agent.service

[Service]
Type=simple
User=yourvpn
Group=yourvpn
WorkingDirectory=/opt/yourvpn/app
EnvironmentFile=/etc/yourvpn/yourvpn.env
ExecStart=/opt/yourvpn/venv/bin/python -m uvicorn yourvpn_api.main:app --host 127.0.0.1 --port 8008 --proxy-headers
Restart=on-failure
RestartSec=3
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/yourvpn /run/yourvpn /var/log/yourvpn

[Install]
WantedBy=multi-user.target
